# ansible/roles/diego/templates/diego.service.j2
[Unit]
Description=Diego service
After=network.target

[Service]
User={{ username }}
ExecStart={{ binary_path }} -interval 5000 -prometheus -prometheus-port {{ port }}
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
